# supervisord.fly.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ── Infrastructure (highest priority, always-on) ──────────────────────────────

[program:postgres]
command=bash -c '
  PG_VER=$(ls /usr/lib/postgresql)
  DATA=/var/lib/postgresql/data
  mkdir -p "$DATA" /var/run/postgresql && chown -R postgres:postgres "$DATA" /var/run/postgresql
  if [ ! -f "$DATA/PG_VERSION" ]; then
    echo "[postgres] Initialising data directory at $DATA"
    su postgres -c "/usr/lib/postgresql/$PG_VER/bin/initdb -D $DATA"
    echo "listen_addresses = '"'"'localhost'"'"'" >> $DATA/postgresql.conf
    su postgres -c "/usr/lib/postgresql/$PG_VER/bin/pg_ctl start -D $DATA -w"
    su postgres -c "psql -c \"CREATE USER veta WITH PASSWORD '"'"'veta'"'"';\""
    su postgres -c "psql -c \"CREATE DATABASE veta OWNER veta;\""
    su postgres -c "/usr/lib/postgresql/$PG_VER/bin/pg_ctl stop -D $DATA"
  fi
  exec su postgres -c "/usr/lib/postgresql/$PG_VER/bin/postgres -D $DATA"'
autostart=true
autorestart=true
priority=1
startsecs=5
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgres.log

[program:redpanda]
command=bash -c 'mkdir -p /var/lib/redpanda/data /etc/redpanda /root/redpanda-wrapper/bin /root/redpanda-wrapper/etc/redpanda && printf "#!/bin/bash\nif [ \"$1\" = \"redpanda\" ]; then shift; fi\nexec /usr/local/lib/redpanda/ld.so --library-path /usr/local/lib/redpanda /usr/local/bin/redpanda \"$@\"\n" > /root/redpanda-wrapper/bin/redpanda && chmod +x /root/redpanda-wrapper/bin/redpanda && /usr/local/bin/rpk redpanda start --install-dir /root/redpanda-wrapper --mode=dev-container --smp=1 --memory=512M --default-log-level=warn --kafka-addr=0.0.0.0:9092 --advertise-kafka-addr=localhost:9092 --advertise-rpc-addr=localhost:33145'
autostart=true
autorestart=true
priority=2
startsecs=10
stdout_logfile=/var/log/supervisor/redpanda.log
stderr_logfile=/var/log/supervisor/redpanda.log

[program:redpanda-check]
command=bash -c 'rm -f /tmp/redpanda-ready; for i in $(seq 1 60); do /usr/local/bin/rpk cluster info --brokers localhost:9092 >/dev/null 2>&1 && touch /tmp/redpanda-ready && exit 0; sleep 3; done; echo "[redpanda-check] timed out waiting for Redpanda"; exit 1'
autostart=true
autorestart=false
priority=3
startsecs=0
exitcodes=0
stdout_logfile=/var/log/supervisor/redpanda-check.log
stderr_logfile=/var/log/supervisor/redpanda-check.log

[program:db-migrate]
command=bash -c 'rm -f /tmp/migration-done; until psql -U veta -d veta -c "SELECT 1" -q 2>/dev/null; do sleep 2; done; cd /app && deno run --allow-all backend/db/migrate.ts && touch /tmp/migration-done'
directory=/app
environment=DATABASE_URL="postgres://veta:veta@localhost:5432/veta"
autostart=true
autorestart=false
priority=3
startsecs=0
exitcodes=0
stdout_logfile=/var/log/supervisor/db-migrate.log
stderr_logfile=/var/log/supervisor/db-migrate.log

; ── Core trading services (always-on) ─────────────────────────────────────────

[program:market-sim]
command=deno run --allow-all /app/backend/src/market-sim/market-sim.ts
directory=/app
autostart=true
autorestart=true
priority=10
startretries=10
stdout_logfile=/var/log/supervisor/market-sim.log
stderr_logfile=/var/log/supervisor/market-sim.log

[program:journal]
command=bash -c 'until [ -f /tmp/migration-done ]; do sleep 2; done; exec deno run --allow-all /app/backend/src/journal/journal-server.ts'
directory=/app
autostart=true
autorestart=true
priority=10
startretries=10
stdout_logfile=/var/log/supervisor/journal.log
stderr_logfile=/var/log/supervisor/journal.log

[program:user-service]
command=bash -c 'until [ -f /tmp/migration-done ]; do sleep 2; done; exec deno run --allow-all /app/backend/src/user-service/user-service.ts'
directory=/app
autostart=true
autorestart=true
priority=10
startretries=10
stdout_logfile=/var/log/supervisor/user-service.log
stderr_logfile=/var/log/supervisor/user-service.log

[program:ems]
command=deno run --allow-all /app/backend/src/ems/ems-server.ts
directory=/app
autostart=true
autorestart=true
priority=10
startretries=10
stdout_logfile=/var/log/supervisor/ems.log
stderr_logfile=/var/log/supervisor/ems.log

[program:oms]
command=deno run --allow-all /app/backend/src/oms/oms-server.ts
directory=/app
autostart=true
autorestart=true
priority=10
startretries=10
stdout_logfile=/var/log/supervisor/oms.log
stderr_logfile=/var/log/supervisor/oms.log

[program:gateway]
command=deno run --allow-all /app/backend/src/gateway/gateway.ts
directory=/app
environment=FRONTEND_DIST="/app/frontend/dist"
autostart=true
autorestart=true
priority=15
startretries=10
stdout_logfile=/var/log/supervisor/gateway.log
stderr_logfile=/var/log/supervisor/gateway.log

; ── Algo strategies ───────────────────────────────────────────────────────────

[group:algos]
programs=algo-trader,twap-algo,pov-algo,vwap-algo,iceberg-algo,sniper-algo,arrival-price-algo,is-algo,momentum-algo

[program:algo-trader]
command=deno run --allow-all /app/backend/src/algo/limit-strategy.ts
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/algo-trader.log
stderr_logfile=/var/log/supervisor/algo-trader.log

[program:twap-algo]
command=deno run --allow-all /app/backend/src/algo/twap-strategy.ts
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/twap-algo.log
stderr_logfile=/var/log/supervisor/twap-algo.log

[program:pov-algo]
command=deno run --allow-all /app/backend/src/algo/pov-strategy.ts
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/pov-algo.log
stderr_logfile=/var/log/supervisor/pov-algo.log

[program:vwap-algo]
command=deno run --allow-all /app/backend/src/algo/vwap-strategy.ts
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/vwap-algo.log
stderr_logfile=/var/log/supervisor/vwap-algo.log

[program:iceberg-algo]
command=deno run --allow-all /app/backend/src/algo/iceberg-strategy.ts
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/iceberg-algo.log
stderr_logfile=/var/log/supervisor/iceberg-algo.log

[program:sniper-algo]
command=deno run --allow-all /app/backend/src/algo/sniper-strategy.ts
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/sniper-algo.log
stderr_logfile=/var/log/supervisor/sniper-algo.log

[program:arrival-price-algo]
command=deno run --allow-all /app/backend/src/algo/arrival-price-strategy.ts
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/arrival-price-algo.log
stderr_logfile=/var/log/supervisor/arrival-price-algo.log

[program:is-algo]
command=deno run --allow-all /app/backend/src/algo/is-strategy.ts
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/is-algo.log
stderr_logfile=/var/log/supervisor/is-algo.log

[program:momentum-algo]
command=deno run --allow-all /app/backend/src/algo/momentum-strategy.ts
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/momentum-algo.log
stderr_logfile=/var/log/supervisor/momentum-algo.log

; ── Microstructure services ───────────────────────────────────────────────────

[group:microstructure]
programs=dark-pool,ccp-service,rfq-service,fix-archive,fix-exchange,fix-gateway

[program:dark-pool]
command=deno run --allow-all /app/backend/src/dark-pool/dark-pool-server.ts
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/dark-pool.log
stderr_logfile=/var/log/supervisor/dark-pool.log

[program:ccp-service]
command=deno run --allow-all /app/backend/src/ccp/ccp-service.ts
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/ccp-service.log
stderr_logfile=/var/log/supervisor/ccp-service.log

[program:rfq-service]
command=deno run --allow-all /app/backend/src/rfq/rfq-service.ts
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/rfq-service.log
stderr_logfile=/var/log/supervisor/rfq-service.log

[program:fix-archive]
command=bash -c 'until [ -f /tmp/migration-done ]; do sleep 2; done; exec deno run --allow-all /app/backend/src/fix/fix-archive.ts'
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/fix-archive.log
stderr_logfile=/var/log/supervisor/fix-archive.log

[program:fix-exchange]
command=deno run --allow-all /app/backend/src/fix/fix-exchange.ts
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/fix-exchange.log
stderr_logfile=/var/log/supervisor/fix-exchange.log

[program:fix-gateway]
command=deno run --allow-all /app/backend/src/fix/fix-gateway.ts
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/fix-gateway.log
stderr_logfile=/var/log/supervisor/fix-gateway.log

; ── Analytics pipeline ────────────────────────────────────────────────────────

[group:analytics]
programs=analytics,market-data-service,market-data-adapters,feature-engine,signal-engine,recommendation-engine,scenario-engine

[program:analytics]
command=bash -c 'until [ -f /tmp/migration-done ]; do sleep 2; done; exec deno run --allow-all /app/backend/src/analytics/analytics-server.ts'
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/analytics.log
stderr_logfile=/var/log/supervisor/analytics.log

[program:market-data-service]
command=deno run --allow-all /app/backend/src/market-data/market-data-service.ts
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/market-data-service.log
stderr_logfile=/var/log/supervisor/market-data-service.log

[program:market-data-adapters]
command=bash -c 'until [ -f /tmp/migration-done ]; do sleep 2; done; exec deno run --allow-all /app/backend/src/market-data-adapters/adapter-server.ts'
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/market-data-adapters.log
stderr_logfile=/var/log/supervisor/market-data-adapters.log

[program:feature-engine]
command=bash -c 'until [ -f /tmp/migration-done ]; do sleep 2; done; exec deno run --allow-all /app/backend/src/feature-engine/feature-engine.ts'
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/feature-engine.log
stderr_logfile=/var/log/supervisor/feature-engine.log

[program:signal-engine]
command=bash -c 'until [ -f /tmp/migration-done ]; do sleep 2; done; exec deno run --allow-all /app/backend/src/signal-engine/signal-engine.ts'
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/signal-engine.log
stderr_logfile=/var/log/supervisor/signal-engine.log

[program:recommendation-engine]
command=bash -c 'until [ -f /tmp/migration-done ]; do sleep 2; done; exec deno run --allow-all /app/backend/src/recommendation-engine/recommendation-server.ts'
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/recommendation-engine.log
stderr_logfile=/var/log/supervisor/recommendation-engine.log

[program:scenario-engine]
command=bash -c 'until [ -f /tmp/migration-done ]; do sleep 2; done; exec deno run --allow-all /app/backend/src/scenario-engine/scenario-server.ts'
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/scenario-engine.log
stderr_logfile=/var/log/supervisor/scenario-engine.log

; ── Aux services ──────────────────────────────────────────────────────────────

[group:aux]
programs=kafka-relay,news-aggregator,llm-advisory-orchestrator,llm-worker

[program:kafka-relay]
command=deno run --allow-all /app/observability/kafka-relay.ts
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/kafka-relay.log
stderr_logfile=/var/log/supervisor/kafka-relay.log

[program:news-aggregator]
command=bash -c 'until [ -f /tmp/migration-done ]; do sleep 2; done; exec deno run --allow-all /app/backend/src/news/news-aggregator.ts'
directory=/app
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/news-aggregator.log
stderr_logfile=/var/log/supervisor/news-aggregator.log

[program:llm-advisory-orchestrator]
command=bash -c 'until [ -f /tmp/migration-done ]; do sleep 2; done; exec deno run --allow-all /app/backend/src/llm-advisory/orchestrator.ts'
directory=/app
environment=LLM_ENABLED="false"
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/llm-advisory.log
stderr_logfile=/var/log/supervisor/llm-advisory.log

[program:llm-worker]
command=deno run --allow-all /app/backend/src/llm-advisory/worker.ts
directory=/app
environment=LLM_ENABLED="false"
autostart=true
autorestart=true
priority=20
startretries=10
stdout_logfile=/var/log/supervisor/llm-worker.log
stderr_logfile=/var/log/supervisor/llm-worker.log
