# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:postgres]
command=bash -c 'mkdir -p /var/run/postgresql && chown postgres:postgres /var/run/postgresql && su postgres -c "/usr/lib/postgresql/$(ls /usr/lib/postgresql)/bin/postgres -D /var/lib/postgresql/$(ls /usr/lib/postgresql)/main -c config_file=/etc/postgresql/$(ls /etc/postgresql)/main/postgresql.conf"'
autostart=true
autorestart=true
priority=2
startsecs=3
stdout_logfile=AUTO
stderr_logfile=AUTO
user=root

[program:redpanda]
command=bash -c 'mkdir -p /var/lib/redpanda/data /etc/redpanda ~/redpanda-wrapper/bin && printf "#!/bin/bash\nif [ \"\$1\" = \"redpanda\" ]; then shift; fi\nexec /usr/local/lib/redpanda/ld.so --library-path /usr/local/lib/redpanda /usr/local/bin/redpanda \"\$@\"\n" > ~/redpanda-wrapper/bin/redpanda && chmod +x ~/redpanda-wrapper/bin/redpanda && /usr/local/bin/rpk redpanda start --install-dir ~/redpanda-wrapper --mode=dev-container --smp=1 --memory=200M --default-log-level=warn --kafka-addr=0.0.0.0:9092 --advertise-kafka-addr=localhost:9092 --advertise-rpc-addr=localhost:33145'
autostart=true
autorestart=true
priority=1
startsecs=8
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:market-sim]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=market-sim deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/market-sim/market-sim.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:ems]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=ems deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/ems/ems-server.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:oms]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=oms deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/oms/oms-server.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:algo-trader]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=limit-algo deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/algo/limit-strategy.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:twap-algo]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=twap-algo deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/algo/twap-strategy.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:pov-algo]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=pov-algo deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/algo/pov-strategy.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:vwap-algo]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=vwap-algo deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/algo/vwap-strategy.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:iceberg-algo]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=iceberg-algo deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/algo/iceberg-strategy.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:sniper-algo]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=sniper-algo deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/algo/sniper-strategy.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:arrival-price-algo]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=arrival-price-algo deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/algo/arrival-price-strategy.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:is-algo]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=is-algo deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/algo/is-strategy.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:momentum-algo]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=momentum-algo deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/algo/momentum-strategy.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

; Deno tcp-proxy permanently holds port 8080 and proxies to Vite on 5173.
; VSCode's port forwarding tunnel never sees 8080 go dark — even when
; Vite restarts, the proxy keeps the TCP listener bound and retries connections.
[program:frontend-proxy]
command=/home/deno/.deno/bin/deno run --allow-net /workspaces/virtual-equities-trading-application/scripts/tcp-proxy.ts 8080 5173
autostart=true
autorestart=true
priority=2
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:frontend]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; deno task dev'
directory=/workspaces/virtual-equities-trading-application/frontend
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:kafka-relay]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=kafka-relay deno run --allow-all /workspaces/virtual-equities-trading-application/observability/kafka-relay.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:user-service]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=user-service deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/user-service/user-service.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:fix-exchange]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=fix-exchange deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/fix/fix-exchange.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:fix-gateway]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=fix-gateway deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/fix/fix-gateway.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:journal]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=journal deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/journal/journal-server.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:fix-archive]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=fix-archive deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/fix/fix-archive.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:gateway]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=gateway deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/gateway/gateway.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:news-aggregator]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=news-aggregator deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/news/news-aggregator.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:analytics]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=analytics deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/analytics/analytics-server.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:market-data-service]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=market-data-service deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/market-data/market-data-service.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:market-data-adapters]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=market-data-adapters deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/market-data-adapters/adapter-server.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:feature-engine]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=feature-engine deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/feature-engine/feature-engine.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:signal-engine]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=signal-engine deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/signal-engine/signal-engine.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:recommendation-engine]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=recommendation-engine deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/recommendation-engine/recommendation-server.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:scenario-engine]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=scenario-engine deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/scenario-engine/scenario-server.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:llm-advisory-orchestrator]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=llm-advisory deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/llm-advisory/orchestrator.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO

[program:llm-worker]
command=bash -c 'source /workspaces/virtual-equities-trading-application/.env 2>/dev/null; OTEL_SERVICE_NAME=llm-worker deno run --allow-all /workspaces/virtual-equities-trading-application/backend/src/llm-advisory/worker.ts'
directory=/workspaces/virtual-equities-trading-application
autostart=false
autorestart=false
stdout_logfile=AUTO
stderr_logfile=AUTO
